# playbooks/blog/roles/srmail/files/srmail-supervisor.conf
[program:srmail]
directory=/usr/local/src/srmail
command=/root/.virtualenvs/srmail/bin/python app/srmail.py config-server.json
