# playbooks/blog/roles/stacosys/files/stacosys-supervisor.conf
[program:stacosys]
directory=/usr/local/src/stacosys
command=/root/.virtualenvs/stacosys/bin/python app/run.py
